# server-logs/brandmen-logs.service
[Unit]
Description=Brandmen logs, updates and device command server
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=-/opt/brandmen-logs/brandmen.env
Environment=PORT=8443
Environment=LOG_DIR=/opt/brandmen-logs/logs
Environment=AGENTS_DIR=/opt/brandmen-logs/agents
Environment=UPDATE_CACHE_DIR=/opt/brandmen-logs/update-cache
Environment=CERT=/opt/brandmen-logs/cert.pem
Environment=KEY=/opt/brandmen-logs/key.pem
Environment=PUBLIC_BASE=https://185.50.203.112
WorkingDirectory=/opt/brandmen-logs
ExecStart=/usr/bin/python3 /opt/brandmen-logs/logserver.py
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
